# systemd/betterforgeprofits-sync.service
[Unit]
Description=BetterForgeProfits Hypixel price sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/betterforgeprofits.env
ExecStart=/usr/bin/docker run --rm --name betterforgeprofits-sync --env-file /etc/betterforgeprofits.env betterforgeprofits-worker:latest sync-prices
Nice=10

[Install]
WantedBy=multi-user.target
